# rpi_client/env.cron
* *  * * *  root /opt/env_monitoring/rpi_client/ds18b20.py
* *  * * *  root /opt/env_monitoring/rpi_client/i2c_arduino.py
* *  * * *  root /opt/env_monitoring/rpi_client/photocell.py
