# rpi_client/env.cron
* *  * * *  root /opt/env_monitoring/rpi_client/18b20.py
