# rpi_client/env.cron
#* *  * * *  root /opt/env_monitoring/rpi_client/ds18b20.py
*/3 *  * * *  root cd /opt/env_monitoring/rpi_client && ./serial_comms.py
*/3 *  * * *  root cd /opt/env_monitoring/rpi_client && ./gps_report.py
*   *  * * *  root cd /opt/env_monitoring/rpi_client && ./hb.py
*/5 *  * * *  root cd /opt/env_monitoring/rpi_client && ./3g.pl
