# rpi_client/env.cron
* *  * * *  root /opt/env_monitoring/rpi_client/ds18b20.py
#* *  * * *  root /opt/env_monitoring/rpi_client/serial_comms.py
